# root/etc/s6-overlay/s6-rc.d/init-mod-universal-tailscale/run
#!/usr/bin/with-contenv bash

## Ubuntu/Debian
if [ -f /usr/bin/apt ]; then
    echo "\
        curl \
        wget \
        iputils-ping \
        dnsutils \
        jq" >> /mod-repo-packages-to-install.list
        # iptables" >> /mod-repo-packages-to-install.list

fi
# Alpine
if [ -f /sbin/apk ]; then
    echo "\
        curl \
        wget \
        jq" >> /mod-repo-packages-to-install.list
        # iptables" >> /mod-repo-packages-to-install.list
fi

# -------------------------------------------
is_true() {
    local var_value=$(eval echo \$$1)
    case "${var_value,,}" in
        1|true|t) return 0 ;;
        *) return 1 ;;
    esac
}
# -------------------------------------------

echo "TAILSCALE_USE_IPTABLES: $TAILSCALE_USE_IPTABLES"
if is_true "TAILSCALE_USE_IPTABLES"; then
    echo "enabling IPTables installation"
    echo "iptables" >> /mod-repo-packages-to-install.list
else
    echo "disabling IPTables installation"
fi